# launchd/com.mitsheth.weekly-cleanup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mitsheth.weekly-cleanup</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/mitsheth/immich-app/backup-scripts/weekly-cleanup.sh</string>
    </array>

    <!-- Run every Sunday at 1:00 AM -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>1</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/mitsheth/immich-app/backup-scripts/logs/launchd-weekly-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/mitsheth/immich-app/backup-scripts/logs/launchd-weekly-stderr.log</string>

    <key>LowPriorityIO</key>
    <true/>
    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
